# recipes-core/busybox/busybox-init/inittab
# /etc/inittab: init configuration.

#<id>:<runlevels>:<action>:<process>

::sysinit:/etc/init.d/rcS
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
